# systemd/unkindness-notify.service
[Unit]
Description=Boot notification service
